# systemd.capf.service
[Unit]
Description=Certificate Authentication Proxy Function
Documentation=https://usecallmanager.nz/certificate-enrollment.html
ConditionPathExists=/var/lib/capf/capf.sqlite3
ConditionPathExists=/var/lib/capf/capf.pem
After=network.target

[Service]
User=capf
ExecStart=/usr/sbin/capfd /var/lib/capfd/capf.sqlite3 -A /var/lib/capf/capf.pem

[Install]
WantedBy=multi-user.target
